# root/etc/services.d/privoxy/run
#!/usr/bin/execlineb -P
privoxy --no-daemon /etc/privoxy/config